# ZBX-8665/misc/init.d/gentoo/zabbix-agentd
#!/sbin/runscript
#
# Zabbix agent start/stop script.
#
# Written by A. Tophofen
# Modified for Zabbix 2.0.0
# May 2012, Zabbix SIA

NAME=zabbix_agentd
PATH=/bin:/usr/bin:/sbin:/usr/sbin:/usr/local/sbin
DAEMON=/usr/local/sbin/${NAME}
DESC="Zabbix agent"
PID=/tmp/$NAME.pid

opts="${opts} reload"

depend() {
	need net
}

start() {
	ebegin "Starting $DESC: $NAME"
	start-stop-daemon --start --pidfile $PID \
		--exec $DAEMON
	eend
}
stop() {
	ebegin "Stopping $DESC: $NAME"
	start-stop-daemon --stop --pidfile $PID \
		--exec $DAEMON
	eend
}
