# ZBX-4203/misc/init.d/ubuntu/zabbix-agent.conf
# zabbix-agent - Start zabbix agent
description     "Zabbix Agent"
author          "S. CANCHON"
start on runlevel [2345]
stop on runlevel [016]
respawn
expect daemon
exec /usr/sbin/zabbix_agentd
